# services/gpsd-client.service
[Unit]
Description=GPSD client service
After=network.target

[Service]
Type=simple
Restart=always
RestartSec=1
ExecStart=gpspipe -w control
StandardOutput=append:/var/log/gps.log
StandardError=append:/var/log/gps.log

[Install]
WantedBy=multi-user.target
