# services/quectel-cm.service
[Unit]
Description=Quectel Connection Manager
Requires=dev-ttyUSB2.device
After=dev-ttyUSB2.device
After=network.target

[Service]
Type=simple
Restart=always
RestartSec=1
ExecStart=/usr/local/bin/quectel-CM -s internet -4
StandardOutput=append:/var/log/quectel-cm.log
StandardError=append:/var/log/quectel-cm.log

[Install]
WantedBy=multi-user.target
